# rebase-helper.service
[Unit]
Description=Rebase-helper: tool which helps you with rebases

[Service]
Type=forking
ExecStart=/usr/bin/rebase-helper-fedmsg.py